# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid
user=root
loglevel=info

; Allow supervisorctl to connect (needed for graceful reload/status checks)
[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

; ---------------------------------------------------------------------------
; php-fpm — must start first so nginx can immediately forward PHP requests
; ---------------------------------------------------------------------------
[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
; stdout/stderr — supervisord 4.x supports /dev/stdout with maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10
startretries=5
startsecs=2

; ---------------------------------------------------------------------------
; nginx — depends on php-fpm being up (priority > php-fpm)
; daemon off; is required when managed by supervisord
; ---------------------------------------------------------------------------
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
startretries=5
startsecs=2
